# files/wallabag.socket
[Unit]
Description=wallabag socket

[Socket]
ListenStream=/run/wallabag.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
