# files/wallabag.service.in
[Unit]
Description=wallabag uwsgi service
Requires=wallabag.socket wallabag-first-run.service
After=wallabag.socket

[Service]
Type=notify
DynamicUser=yes
User=wallabag

TemporaryFileSystem=/var/lib
StateDirectory=wallabag
Environment=WALLABAG_DATA=/var/lib/wallabag
EnvironmentFile=/etc/default/wallabag.conf

ExecStart=@execStart@
Restart=always

[Install]
WantedBy=multi-user.target
Also=wallabag.socket
